# sys-cluster/seaweedfs-bin/files/seaweedfs-volume.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

supervisor=supervise-daemon

name=weed
command="/opt/seaweedfs-bin/weed"
command_args="${WEED_EXEC} ${WEED_OPTS}"

pidfile="/var/run/${RC_SVCNAME}.pid"
respawn_delay=5

depend() {
	need localmount
	need net
}
